# system/nerdy_things_ai_transcription.service
[Unit]
Description=Run transcoding queue
After=multi-user.target

[Service]
ExecStart=/usr/bin/python3 /home/__USER__/raspberry-pi-audio-spy/python/daemon_ai.py 0
WorkingDirectory=/home/__USER__/raspberry-pi-audio-spy/python/
StandardOutput=file:/home/__USER__/raspberry-pi-audio-spy/system/logs/daemon_ai.log
StandardError=file:/home/__USER__/raspberry-pi-audio-spy/system/logs/daemon_ai.log
Restart=always
User=__USER__
RemainAfterExit=True

[Install]
WantedBy=multi-user.target
